# deploy/secretary-bot.service
[Unit]
Description=KGB_Bot — Telegram knowledge-base bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/KGB_Bot
# Secrets are read by the app from ./.env (pydantic-settings). The file must be
# chmod 600 and owned by the service user. WorkingDirectory must be the repo root.
ExecStart=/home/YOUR_USER/KGB_Bot/.venv/bin/python -m secretary_bot
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/YOUR_USER/KGB_Bot/data /home/YOUR_USER/KGB_Bot/.env

[Install]
WantedBy=multi-user.target
