# deploy/secretary-bot.user.service
# User-level systemd unit (no sudo needed; requires `loginctl enable-linger`).
# Install: ~/.config/systemd/user/secretary-bot.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now secretary-bot
[Unit]
Description=KGB_Bot — Telegram knowledge-base bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/KGB_Bot
ExecStart=%h/KGB_Bot/.venv/bin/python -m secretary_bot
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=default.target
